# AI-projects/ai-system-health-dashboard/agents/debian/ashd-agent.service
[Unit]
Description=System Trace Monitoring Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/system-trace-agent
ExecStart=/usr/bin/python3 /opt/system-trace-agent/system-trace_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
